# deploy/operatoros.service
[Unit]
Description=OperatorOS Next.js app
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/operatoros
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
Environment=OPERATOROS_DATABASE_PATH=/var/lib/operatoros/operatoros-db.sqlite
Environment=OPERATOROS_BASE_URL=https://operatoros.example.com
# Prefer OPERATOROS_USERS_JSON via a systemd drop-in or EnvironmentFile.
# Legacy OWNER_PASSWORD fallback is still supported.
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=5
User=operatoros
Group=operatoros

[Install]
WantedBy=multi-user.target
